# docker/image/assets/s6-rc.d/sway/run
#!/usr/bin/with-contenv bash
exec /deployery/start-sway.sh
