# docker/image/assets/s6-rc.d/code-server/run
#!/usr/bin/with-contenv bash
exec /deployery/start-code-server.sh
