# docker/image/assets/s6-rc.d/init-deployery/up
#!/usr/bin/with-contenv bash
exec /entrypoint.sh
